# filescanner.service
[Unit]
Description=dosya tarayici

[Service]
Type=simple
ExecStart=/home/aleyna/dosyatarayici.sh

[Install]
WantedBy=multi-user.target
